# misc/package-scripts/preremove
#!/bin/bash

# remove systemd unit
systemctl disable netconfig.service
rm /usr/lib/systemd/system/netconfig.service

